# scripts/deliver/com.user.agent-team.order-check.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.agent-team.order-check</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>{{REPO_DIR}}/scripts/deliver/auto_check_orders.sh</string>
    </array>

    <!-- 1日5回：8:30, 12:00, 15:00, 18:00, 21:00 に自動実行 -->
    <!-- PCが寝てた場合、起動時にキャッチアップ実行 -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>8</integer>
            <key>Minute</key>
            <integer>30</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>12</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>15</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>18</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>21</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>{{REPO_DIR}}/logs/order_check.out.log</string>

    <key>StandardErrorPath</key>
    <string>{{REPO_DIR}}/logs/order_check.err.log</string>

    <key>WorkingDirectory</key>
    <string>{{REPO_DIR}}</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
